# datanode/supervisord.conf
[supervisord]
nodaemon=true

[program:datanode]
command=hdfs datanode
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
