# historyserver/supervisord.conf
[supervisord]
nodaemon=true

[program:historyserver]
command=mapred historyserver
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

