# namenode/supervisord.conf
[supervisord]
nodaemon=true

[program:namenode]
command=hdfs namenode
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

