# nodemanager/supervisord.conf
[supervisord]
nodaemon=true

[program:nodemanager]
command=yarn nodemanager
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

